# scripts/com.gia.learning.daughter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gia.learning</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/gialee/NSS_Learning/NSS_Word_Master/launcher/gia-daughter-launch.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/gialee/NSS_Learning/NSS_Word_Master</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/gialee</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>15</integer>
  <key>StandardOutPath</key>
  <string>/Users/gialee/NSS_Learning/logs/gia-stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/gialee/NSS_Learning/logs/gia-stderr.log</string>
  <key>Nice</key>
  <integer>5</integer>
</dict>
</plist>
